# apps/slack-bot/deploy/ailms-bot.service
[Unit]
Description=AILMS Slack Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/ailms

# Load all env vars from the root .env file
EnvironmentFile=/opt/ailms/.env

# Run via pnpm so workspace module resolution works correctly
ExecStart=/usr/local/bin/pnpm --filter slack-bot run start

# Restart automatically on crash, wait 10s between attempts
Restart=always
RestartSec=10

StandardOutput=journal
StandardError=journal
SyslogIdentifier=ailms-bot

[Install]
WantedBy=multi-user.target
